# ansible/roles/dci-api/templates/dci.service
[Unit]
Description=Distributed CI system

[Service]
PIDFile=/var/run/dci.pid
ExecStart={{path_to_dci}}/runtestserver.py
